# config/deploy.service
# ABOUTME: Systemd service file for running Deploy in production.
# ABOUTME: Copy to /etc/systemd/system/ and enable with systemctl.

[Unit]
Description=Deploy Server
Documentation=https://github.com/keithk/deploy
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/deploy
ExecStart=/home/deploy/.bun/bin/bun run start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Load environment variables from .env (created by `deploy setup`)
EnvironmentFile=/home/deploy/deploy/.env

# Environment variables (also set in .env; these are fallbacks)
Environment=NODE_ENV=production
Environment=BUN_INSTALL=/home/deploy/.bun
Environment=PATH=/home/deploy/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/deploy/deploy/data /home/deploy/deploy/logs /var/deploy/sites /var/deploy/data /var/log/caddy
PrivateTmp=true

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
